# debian/philipplay.service
[Unit]
Description=Philipplay Audio Player

[Service]
Type=simple
ExecStart=/usr/bin/philipplay

[Install]
WantedBy=multi-user.target
